# Linux/logwatcher-0.3.1/watcher.cron
SHELL=/bin/bash
*/1 * * * * root /usr/local/sbin/watcher > /dev/null 2>&1
