# Security/dfp-0.3/dfp.cron
SHELL=/bin/bash
*/1 * * * * root /usr/local/sbin/dfp > /dev/null 2>&1
